# deploy/gm-backend.service
[Unit]
Description=GM Stream Control — FastAPI Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/gm-stream-control/backend
ExecStart=/opt/gm-stream-control/backend/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
EnvironmentFile=/opt/gm-stream-control/backend/.env

[Install]
WantedBy=multi-user.target
